# src/etc/skel/.config/systemd/user/default.target.wants/bbqlinux-mate-settings.service
[Unit]
Description=Apply BBQLinux specific MATE settings

[Service]
Type=simple
ExecStart=/usr/bin/bbqlinux-mate-settings.sh

[Install]
WantedBy=default.target
